# nebula-watcher.path
[Path]
PathModified=/etc/nebula/ansible_date

[Install]
WantedBy=multi-user.target